# wn-printerhub.service
# systemd service file for WN-PrinterHub
# Copy to /etc/systemd/system/wn-printerhub.service

[Unit]
Description=WN-PrinterHub - Local LAN printer connector service
Documentation=https://github.com/whiteneuron/WN-PrinterHub
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=printerhub
Group=printerhub
WorkingDirectory=/opt/wn-printerhub
Environment=PATH=/opt/wn-printerhub/.venv/bin
EnvironmentFile=/opt/wn-printerhub/.env
ExecStart=/usr/bin/python3 /opt/wn-printerhub/production.py
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/wn-printerhub /var/log
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wn-printerhub

[Install]
WantedBy=multi-user.target